# supervisord.conf
[supervisord]
nodaemon=true

[program:entry_point]
command=/opt/bin/entry_point.sh

[program:sikulixapi]
command=/usr/bin/java -jar /root/SikuliX/sikulixapi.jar -pythonserver
